# data/system-updates.service.in
[Unit]
Description=System updates server

[Service]
ExecStart=@libexecdir@/dev.rlxos.updates
BusName=dev.rlxos.updates
Restart=on-failure

[Install]
WantedBy=multi-user.target